# image/service/sshd/run
#!/bin/sh
set -e
exec /usr/sbin/sshd -D
